# server/raedsync-backup.timer
[Unit]
Description=Run Raedworkouts sync backups hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
